# data/idea/debian/postinst
#!/bin/bash

. /usr/share/debconf/confmodule

echo "Reloading sysctl settings for inotify settings."
invoke-rc.d procps start

db_input high idea-EDITION/java_hint || true
db_go
